# ansible/service_repository/harbor.service
[Unit]
Description=Harbor Service
After=network.target docker.service
[Service]
Type=simple
WorkingDirectory=/opt/harbor
ExecStart=/usr/bin/docker-compose up
ExecStop=/usr/bin/docker-compose down
RemainAfterExit=yes
[Install]
WantedBy=multi-user.target
